# services.d/http/run
#!/usr/bin/execlineb -P

/usr/local/bin/php -S 0.0.0.0:80 -t /satis/data
