# services.d/satis-build/run
#!/usr/bin/with-contenv sh

/root/.composer/vendor/bin/satis build /satis/satis.json /satis/data
sleep 300 # 5 minutes
