# services.d/satis-build/run
#!/usr/bin/with-contenv sh

/satis/bin/satis build /satis.json /data
sleep 300 # 5 minutes
